# systemd_services/grafana_daily_backup.timer
[Unit]
Description=grafana_daily_backup timer

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target